# backup/crontab
# Backup every Monday morning
# min   hour    day     month   weekday command
0     0       *       *       1       /skatliste-backup.sh
